# SigmaStart_Rootfs/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:-/bin/sh
#console::respawn:/sbin/getty -L console 115200 vt100
